# install/autostart/lemca_autostart.desktop
[Desktop Entry]
Type=Application
Name=agtf
Exec=sh -c "~/lemca/build/lemca"
Terminal=true
Hidden=false
X-GNOME-Autostart-enabled=true
